# dev-db/monetdb-sql/files/monetdb.init-5.2.0
#!/sbin/runscript
# Copyright 1999-2007 Gentoo Foundation; Distributed under the GPL v2
# $Header: $

depend() {
	use net
}

start() {
	ebegin "Starting MonetDB/SQL"
	# -b can go away next release, hopefully, see:
	# http://sourceforge.net/tracker/index.php?func=detail&aid=1825922&group_id=56967&atid=482471
	start-stop-daemon --start \
		-c $MDB_USER:$MDB_GROUP \
		-b \
		-m -p /var/run/merovingian \
		--exec /usr/sbin/merovingian
	eend $?
}

stop() {
	ebegin "Stopping MonetDB/SQL"
	start-stop-daemon --stop \
		-c $MDB_USER:$MDB_GROUP \
		-p /var/run/merovingian \
		--exec /usr/sbin/merovingian
	eend $?
}
